# docker-kknockfile/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true

[program:apache2]
command=/usr/sbin/httpd -D FOREGROUND
autostart=true
autorestart=true
